# nginx-php-supervisor/etc/supervisor/conf.d/nginx-php-fmp.conf
[supervisord]
nodaemon=true

[program:php-fpm7.0]
command=/bin/bash -c "exec /usr/sbin/php-fpm7.0 --nodaemonize"

[program:nginx]
command=/bin/bash -c "exec /usr/sbin/nginx -g \"daemon off;\""
