# nginx-ping-supervisor/etc/supervisor/conf.d/nginx.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/bin/bash -c "exec /usr/sbin/nginx -g \"daemon off;\""
